# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
; Run migrations first, then start the API server
command=sh -c "alembic upgrade head && uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2"
directory=/app/backend
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app/backend"
; nginx waits for backend to be healthy before getting traffic
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
; Give backend 5s to bind before nginx starts accepting connections
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
